# roles/gui_basic/tasks/dotfiles/backup/systemdTimer/configUpload.timer
[Unit]
Description=configUpload
Wants=configUpload.service

[Timer]
Unit=configUpload.service
OnCalendar=daily

[Install]
WantedBy=timers.target
